# systemd/mail2slack@.path
[Path]
PathModified=/var/mail/%i

[Install]